# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:health]
command=python -m backend.services.health.main
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/health.err.log
stdout_logfile=/var/log/supervisor/health.out.log

default_startsecs=5

[program:api_gateway]
command=python start_api_gateway.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/api_gateway.err.log
stdout_logfile=/var/log/supervisor/api_gateway.out.log
environment=DOCKER_COMPOSE=true
depends_on=health

default_startsecs=5

[program:price_feeds]
command=python -m backend.services.price_feeds.main
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/price_feeds.err.log
stdout_logfile=/var/log/supervisor/price_feeds.out.log
environment=DOCKER_COMPOSE=true,API_GATEWAY_URL=http://localhost:5000
depends_on=api_gateway

default_startsecs=5

[program:arbitrage]
command=python -m backend.services.arbitrage.main
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/arbitrage.err.log
stdout_logfile=/var/log/supervisor/arbitrage.out.log
depends_on=price_feeds

default_startsecs=5

[program:cep_engine]
command=python -m backend.services.cep_engine.main
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/cep_engine.err.log
stdout_logfile=/var/log/supervisor/cep_engine.out.log
depends_on=arbitrage

default_startsecs=5

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log
depends_on=api_gateway

default_startsecs=5 